# roles/studentenportal/files/studentenportal-sqlbackup.timer
[Unit]
Description=Run SQL backup daily

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
