# roles/studentenportal/files/studentenportal-dehydrated.timer
[Unit]
Description=Run dehydrated update daily

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target