# roles/backups/files/duply.timer
[Unit]
Description=Run duply weekly

[Timer]
OnCalendar=weekly

[Install]
WantedBy=timers.target
